# polkit-qt-example.service
[Unit]
Description=polkit-qt-example-daemon

[Service]
Type=dbus
BusName=com.polkit.qt.example.dbus
